# systemd/ipv6-frag-escape-tracker-update.service
[Unit]
Description=Update IPV6_FRAG_ESCAPE tracker via Claude Code

[Service]
Type=oneshot
WorkingDirectory=%h/src/auto-update/ipv6_frag_escape
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/src/auto-update/ipv6_frag_escape/scripts/auto-update
StandardOutput=journal
StandardError=journal
